# aws/mongo/etc/cron.d/spinpunch
# run daily tasks at 0700 UTC time (usually lowest load)
0 7 * * * ec2-user run-parts /etc/cron.spinpunch.daily
# 30 * * * * ec2-user run-parts /etc/cron.spinpunch.hourly

# send memory statistics to CloudWatch every minute
* * * * * nobody /usr/local/bin/ec2-send-memory-metrics.py
