# aws/gamemaster/etc/cron.d/spinpunch
0,5,10,15,20,25,30,35,40,45,50,55 * * * * ec2-user run-parts /etc/cron.spinpunch.rapid
0 * * * * ec2-user run-parts /etc/cron.spinpunch.hourly

