# aws/gamemaster/etc/cron.daily/spinpunch-backup-corp-data.sh
#!/bin/sh

/home/ec2-user/ransomerrift/aws/backup-corp-data.sh > /dev/null

